# site-template/systemd/recsync.service.in
[Unit]
Description=RecSync Service
Documentation=@DOCURL1@
Documentation=@DOCURL2@
After=network.target @SYSTEMD_SERVICES@
Requires=@SYSTEMD_SERVICES@

[Service]
#
# Twistd also has the options, but we use the systemd one instead.
#
User=@USERID@
Group=@GROUPID@

# Everything in a local, so we have to tell where PYTHONPATH is
Environment="PYTHONPATH=@RS_INSTALL_LOCATION@"

# Due to User/Group permission, we cannot generate the recsync.log
# in /var/log/ folder. So we put in into the recsync installation
# location instead. In the same way, a PID file also cannot be
# generated. However, we don't need it, since systemd service will
# take care of this. 
# 
ExecStart=@TWISTD_PATH@/twistd --logfile=@RS_INSTALL_LOCATION@/recsync.log \
			       --rundir=@RS_INSTALL_LOCATION@ \
			       --pidfile= \
			       --reactor=poll \
			       --nodaemon recceiver -f @RS_INSTALL_LOCATION@/@RECSYNC_CONF@

Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
Alias=rs.service
#
