# site-template/rs.service.in
[Unit]
Description=RecSync Service
Documentation=https://github.com/ChannelFinder/recsync/tree/master/server
Documentation=https://github.com/ChannelFinder/pyCFClient
After=network.target channelfinder.service
Wants=channelfinder.service

[Service]

#
# Twistd also has the options, but we use the systemd one instead.
#
User=_USERID_
Group=_GROUPID_

# Everything in a local, so we have to tell where PYTHONPATH is
Environment="PYTHONPATH=_PYTHONPATH_"

# Due to User/Group permission, we cannot generate the recsync.log
# in /var/log/ folder. So we put in into the recsync installation
# location instead. In the same way, a PID file also cannot be
# generated. However, we don't need it, since systemd service will
# take care of this. 
# 
ExecStart=_TWISDPATH_/twistd   --logfile=_PYTHONPATH_/recsync.log \
			       --rundir _PYTHONPATH_ \
			       --pidfile= \
			       --reactor=poll \
			       --nodaemon recceiver -f _PYTHONPATH_/_CONFIGFILE_

Restart=always

[Install]
WantedBy=multi-user.target
Alias=rs.service
#
